# 네트워크/리눅스/crontab
일정한 시간 반복, 특정 시간에 맞게 쉘 스크립트를 실행
확인 및 설정: crontab -e

예시 : 매일 새벽 1시에 해당 경로에 있는 쉘 스크립트 실행
0 1 * * * /쉘스크립트 경로
